# example-systemd-unit/sa-learn-nicely.service
[Unit]
Description=sa-learn one mail at a time

[Service]
Type=simple
ExecStart=/dis/mail--sa-learn-nicely/bin/sa-learn-nicely

